# openvpn-client@tb.service.d/10-netns.conf
[Unit]
Requires=netns@tb.service
After=netns@tb.service
JoinsNamespaceOf=netns@tb.service

Requires=netns-attach@tb.service
After=netns-attach@tb.service

[Service]
PrivateNetwork=true

